# src/systemd/.config/systemd/user/slideshow-watcher.path
[Path]
PathModified=/home/pi/Pictures/ImageLib/slideshow-all.jar-updated

[Install]
WantedBy=multi-user.target
