# src/deployment/.config/systemd/user/slideshow-watcher.path
[Path]
PathModified=/home/pi/slideshow-all.jar-updated

[Install]
WantedBy=multi-user.target
